# relay/com.wechat.relay.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wechat.relay</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mint/.local/bin/uv</string>
        <string>run</string>
        <string>--with</string>
        <string>pypsrp</string>
        <string>--with</string>
        <string>sherpa-onnx</string>
        <string>python</string>
        <string>-u</string>
        <string>relay_watch.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mint/workspace/tools/wechat-copilot/relay</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/mint/workspace/tools/wechat-copilot/relay/relay_watch.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mint/workspace/tools/wechat-copilot/relay/relay_watch.log</string>
</dict>
</plist>
